# nginx/configs/ini/crond.ini
[program:crond]
command=crond -f -l 0
autostart=true
autorestart=true
