# nginx/configs/ini/crond.ini
[program:crond]
command=crond -f -l 0
autostart=true
autorestart=true
stderr_logfile=/var/log/crond.err.log
stdout_logfile=/var/log/crond.out.log

