# nginx/configs/ini/nginx.ini
[program:nginx]
command=/scripts/start-nginx.sh
autostart=true
autorestart=true
startretries=2
