# nginx/configs/ini/nginx.ini
[program:nginx]
command=/scripts/start-nginx.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/nginx.err.log
stdout_logfile=/var/log/nginx.out.log
