# nginx/configs/ini/nginx.ini
[program:nginx]
command=/scripts/start-nginx.sh
autostart=true
autorestart=true
startretries=2
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
