# nginx/configs/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/scripts/start-nginx.sh
autostart=true
autorestart=true

[program:crond]
command=crond -f -l 0
autostart=true
autorestart=true
